# skunkyart.example.openrc
#!/sbin/openrc-run
name="SkunkyArt"
description="Privacy frontend for deviantart.com"
supervisor=supervise-daemon
command=<path_to_skunkyart>
command_args="-c <path_to_config>"
directory="<path_to_dir_with_skunkyart>"

depend() {
	need net
}