# services/skunkyart.example.openrc
#!/sbin/openrc-run
supervisor=supervise-daemon
user=skunkyart:skunkyart
name=SkunkyArt
directory=/opt/skunkyart
command=$directory/bin/skunkyart
description="Privacy-oriented frontend for DeviantArt"
error_logger=logger
output_logger=logger
no_new_privs=true
umask=0077

# if you use old version openrc, uncomment function lower
# start_pre() {
# 	cd $directory
# }

depend() {
    need net localmount bootmisc
}